# config/crontab
# NTP time sync
@hourly ntpd -d -n -q -p pool.ntp.org >/dev/null 2>&1

# Restic backup
@hourly /restic.sh backup
