# config/crontab
# Restic backup
@hourly /restic.sh backup
